# doc/sample_upstart_script/assembl.conf
description     "Assembl supervisord upstart script"

# This is the line you have to modify for your environment
env ASSEMBL_PATH="/var/www/assembl"

start on runlevel [345]
stop on runlevel [!345]

setuid www-data
setgid www-data

respawn

script
  chdir $ASSEMBL_PATH
  # Yes, we really do need this separate activate line so supervisord
  # cab execute shell scripts from the same venv
  . venv/bin/activate
  exec supervisord --nodaemon
end script

post-start script
  chdir $ASSEMBL_PATH
  . venv/bin/activate
  #Nowhere near perfect, but we can't use supervisorctl's return code
  sleep 1
  exec supervisorctl start prod:uwsgi
end script

